# robot/services/volume-max.service
[Unit]
Description=Set HifiBerry as default sink and max volume (PipeWire)
After=sound.target graphical-session.target
Wants=sound.target

[Service]
Type=oneshot
User=pi
Group=audio
SupplementaryGroups=audio
# PipeWire/WirePlumber, not ALSA 'Master' (which doesn't exist on the HifiBerry).
ExecStart=/home/pi/pupperv3/robot/services/set_audio_default.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="HOME=/home/pi"
Environment="XDG_RUNTIME_DIR=/run/user/1000"

[Install]
WantedBy=multi-user.target
