# pupper-rs/pupper-gui.service
[Unit]
Description=Pupper GUI Service
After=graphical-session.target sound.target network.target
Wants=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/home/pi/pupperv3-monorepo/pupper-rs
ExecStart=/home/pi/pupperv3-monorepo/pupper-rs/run_gui.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment="HOME=/home/pi"
Environment="USER=pi"
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/pi/.cargo/bin"
Environment="DISPLAY=:0"
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="XDG_SESSION_TYPE=wayland"
Environment="ROS_LOCALHOST_ONLY=1"

User=pi
Group=pi
SupplementaryGroups=i2c video audio render gpio spi

PrivateDevices=no
# Disable some systemd hardening for hardware access
NoNewPrivileges=no
ProtectKernelModules=no

[Install]
WantedBy=graphical.target